# sherpa-onnx-tts-stt/rootfs/etc/s6-overlay/s6-rc.d/sherpa-onnx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# ==============================================================================
declare config
LANGUAGE=$(echo `bashio::config 'language'`|sed "s/null/$LANGUAGE/g")
SPEED=$(echo `bashio::config 'speed'`|sed "s/null/$SPEED/g")
STT_MODEL=$(echo `bashio::config 'stt_model'`|sed "s/null/$STT_MODEL/g")
STT_USE_INT8_ONNX_MODEL=$(echo `bashio::config 'stt_use_int8_onnx_model'`|sed "s/null/$STT_USE_INT8_ONNX_MODEL/g")
STT_THREAD_NUM=$(echo `bashio::config 'stt_thread_num'`|sed "s/null/$STT_THREAD_NUM/g")
TTS_MODEL=$(echo `bashio::config 'tts_model'`|sed "s/null/$TTS_MODEL/g")
TTS_THREAD_NUM=$(echo `bashio::config 'tts_thread_num'`|sed "s/null/$TTS_THREAD_NUM/g")
TTS_SPEAKER_SID=$(echo `bashio::config 'tts_speaker_sid'`|sed "s/null/$TTS_SPEAKER_SID/g")
DEBUG=$(echo `bashio::config 'debug'`|sed "s/null/$DEBUG/g")
exec python3 /app/run.py \
    --language "$LANGUAGE" \
    --speed "$SPEED" \
    --stt_model "$STT_MODEL" \
    --stt_use_int8_onnx_model "$STT_USE_INT8_ONNX_MODEL" \
    --stt_thread_num "$STT_THREAD_NUM" \
    --tts_model "$TTS_MODEL" \
    --tts_thread_num "$TTS_THREAD_NUM" \
    --tts_speaker_sid "$TTS_SPEAKER_SID" \
    --debug "$DEBUG" 
