# sherpa-onnx-tts-stt/rootfs/etc/s6-overlay/s6-rc.d/sherpa-onnx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# ==============================================================================
declare config

exec python3 /app/run.py \
    --language "$(bashio::config 'language')" \
    --speed "$(bashio::config 'speed')" \
    --stt_model: "$(bashio::config 'stt_model')" \
    --stt_use_int8_onnx_model: "$(bashio::config 'stt_use_int8_onnx_model')" \
    --stt_thread_num: "$(bashio::config 'stt_thread_num')" \
    --tts_model: "$(bashio::config 'tts_model')" \
    --tts_thread_num: "$(bashio::config 'tts_thread_num')" \
    --tts_speaker_sid: "$(bashio::config 'tts_speaker_sid')" \
    --debug: "$(bashio::config 'debug')" 
