# sherpa-onnx-tts-stt/rootfs/etc/s6-overlay/s6-rc.d/sherpa-onnx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# ==============================================================================
declare config

python3 /app/run.py
