# sherpa-onnx-tts-stt/rootfs/etc/s6-overlay/s6-rc.d/sherpa-onnx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# ==============================================================================
declare config

exec python3 /app/run.py \
    --language "${$(bashio::config 'language'):=$LANGUAGE}" \
    --speed "${$(bashio::config 'speed'):=$SPEED}" \
    --stt_model: "${$(bashio::config 'stt_model'):=$STT_MODEL}" \
    --stt_use_int8_onnx_model: "${$(bashio::config 'stt_use_int8_onnx_model'):=$STT_USE_INT8_ONNX_MODEL}" \
    --stt_thread_num: "${$(bashio::config 'stt_thread_num'):=$STT_THREAD_NUM}" \
    --tts_model: "${$(bashio::config 'tts_model'):=$TTS_MODEL}" \
    --tts_thread_num: "${$(bashio::config 'tts_thread_num'):=$TTS_THREAD_NUM}" \
    --tts_speaker_sid: "${$(bashio::config 'tts_speaker_sid'):=$TTS_SPEAKER_SID}" \
    --debug: "${$(bashio::config 'debug'):=$DEBUG}" 
