# sherpa-onnx-tts-stt/rootfs/etc/s6-overlay/s6-rc.d/sherpa-onnx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# ==============================================================================
declare config
LANGUAGE=$(echo `bashio::config 'language'`|sed "s/null/$LANGUAGE/g")
SPEED=$(echo `bashio::config 'speed'`|sed "s/null/$SPEED/g")
STT_MODEL=$(echo `bashio::config 'stt_model'`|sed "s/null/$STT_MODEL/g")
STT_USE_INT8_ONNX_MODEL=$(echo `bashio::config 'stt_use_int8_onnx_model'`|sed "s/null/$STT_USE_INT8_ONNX_MODEL/g")
STT_THREAD_NUM=$(echo `bashio::config 'stt_thread_num'`|sed "s/null/$STT_THREAD_NUM/g")
TTS_MODEL=$(echo `bashio::config 'tts_model'`|sed "s/null/$TTS_MODEL/g")
TTS_THREAD_NUM=$(echo `bashio::config 'tts_thread_num'`|sed "s/null/$TTS_THREAD_NUM/g")
TTS_SPEAKER_SID=$(echo `bashio::config 'tts_speaker_sid'`|sed "s/null/$TTS_SPEAKER_SID/g")
DEBUG=$(echo `bashio::config 'debug'`|sed "s/null/$DEBUG/g")

if [ "$(bashio::config 'custom_stt_model')" != 'null' ]; then
    CUSTOM_STT_MODEL=$(bashio::config 'custom_stt_model')
fi

if [ "$(bashio::config 'custom_stt_model_eval')" != 'null' ]; then
    CUSTOM_STT_MODEL_EVAL=$(bashio::config 'custom_stt_model_eval')
fi

if [ "$(bashio::config 'custom_tts_model')" != 'null' ]; then
    CUSTOM_TTS_MODEL=$(bashio::config 'custom_tts_model')
fi

if [ "$(bashio::config 'custom_tts_model_eval')" != 'null' ]; then
    CUSTOM_TTS_MODEL_EVAL=$(bashio::config 'custom_tts_model_eval')
fi

echo "$LANGUAGE" "$SPEED"  "$STT_MODEL"  "$STT_USE_INT8_ONNX_MODEL"  "$STT_THREAD_NUM"  "$TTS_MODEL"  "$TTS_THREAD_NUM"  "$TTS_SPEAKER_SID"  "$DEBUG" "$CUSTOM_STT_MODEL" "$CUSTOM_STT_MODEL_EVAL" "$CUSTOM_TTS_MODEL" "$CUSTOM_TTS_MODEL_EVAL" 
exec python3 /app/run.py \
    --language "$LANGUAGE" \
    --speed "$SPEED" \
    --stt_model "$STT_MODEL" \
    --stt_use_int8_onnx_model "$STT_USE_INT8_ONNX_MODEL" \
    --stt_thread_num "$STT_THREAD_NUM" \
    --tts_model "$TTS_MODEL" \
    --tts_thread_num "$TTS_THREAD_NUM" \
    --tts_speaker_sid "$TTS_SPEAKER_SID" \
    --debug "$DEBUG" \
    --custom_stt_model "$CUSTOM_STT_MODEL" \
    --custom_stt_model_eval "$CUSTOM_STT_MODEL_EVAL" \
    --custom_tts_model "$CUSTOM_TTS_MODEL" \
    --custom_tts_model_eval "$CUSTOM_TTS_MODEL_EVAL"
